# deploy/smart_home-refresh.service
[Unit]
Description=smart_home daily day-ahead plan refresh (ENTSO-E)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# the retry loop can run until 17:00, so disable the default start timeout
TimeoutStartSec=infinity
User=solarpi
WorkingDirectory=/home/solarpi/smart_home
EnvironmentFile=/etc/smart_home.env
ExecStart=/home/solarpi/smart_home/.venv/bin/python -m smart_home.schedule refresh --wait-tomorrow
