# deploy/smart_home-refresh.timer
[Unit]
Description=Refresh the day-ahead plan daily (start 12:15; retry until tomorrow's prices publish)

[Timer]
# Start at 12:15 — before BE day-ahead prices publish (~12:45-13:00 CET). The service then
# retries (5 min for the first hour, then hourly) until tomorrow's prices are cached or 17:00.
OnCalendar=*-*-* 12:15:00
Persistent=true

[Install]
WantedBy=timers.target
