# deploy/smart_home-refresh.timer
[Unit]
Description=Refresh the day-ahead plan daily (tomorrow's prices publish ~13:00 CET)

[Timer]
OnCalendar=*-*-* 16:00:00
# Spread the fire time over 10 min so we don't hammer ENTSO-E at the top-of-hour
# spike (everyone polls at HH:00:00 right as day-ahead data settles).
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
